# utils/systemd-units/sshd@.service
[Unit]
Description=SSH Per-Connection Server

[Service]
ExecStart=/usr/bin/sshd -i
StandardInput=socket
StandardError=syslog
